# local/uu-uuid-persist/uuplugin-uuid-sync.path
[Unit]
Description=Watch UU Plugin UUID file

[Path]
PathModified=/tmp/uu/.uuplugin_uuid
