# local/uu-uuid-persist/uuplugin-uuid-sync.service
[Unit]
Description=Backup UU Plugin UUID

[Service]
Type=oneshot
ExecStart=/usr/local/bin/uuplugin-uuid-helper.sh backup
